# com.example.morning_playlist_update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.morning_playlist_update</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/path/to/Spotify-music-list/morning_playlist_update.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>50</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/morning_playlist_update.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/morning_playlist_update.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <!-- put CLIENT_ID, CLIENT_SECRET, REDIRECT_URI here if not using .env -->
    </dict>
</dict>
</plist>
